# macos/LaunchAgents/com.zhiwen.obsidian.bisync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- 必须的唯一标识 -->
    <key>Label</key>
    <string>com.zhiwen.obsidian.bisync</string>

    <!-- 程序 + 参数：全部写绝对路径 -->
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/zhiwen.wang/GithubRepos/dotfiles/macos/sync_dropbox_icloud_obsidian.sh</string>
    </array>

    <!-- 每 30 分钟跑一次 -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- 登录后立即跑一次 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 日志（先确保目录存在） -->
    <key>StandardOutPath</key>
    <string>/Users/zhiwen.wang/Library/Logs/obsidian_bisync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zhiwen.wang/Library/Logs/obsidian_bisync.err</string>
  </dict>
</plist>